# root/etc/uci-defaults/10-luci-app-twofa
#!/bin/sh
# Initialize twofa config and grant rpcd ACL for root.
# Never run uci against /etc/config/twofa until the file parses cleanly.

name="luci-app-twofa"
changed_twofa=0
changed_rpcd=0

write_twofa_config() {
	cat > /etc/config/twofa <<'EOF'
config twofa 'global'
	option enabled '0'
	option issuer 'OpenWrt'
EOF
}

generate_secret() {
	# 16-char base32 secret for TOTP.
	#
	# IMPORTANT: openssl rand has -hex and -base64 but NOT -base32, so the
	# previous `openssl rand -base32 20` call returned an empty string and
	# silently committed `secret=""` into UCI - that's how router and
	# authenticator could end up holding *different* secrets even when the
	# user never explicitly regenerated.
	#
	# Pull random bytes from /dev/urandom and let tr filter them to the
	# base32 alphabet directly. tr is in busybox on every OpenWrt build.
	local s
	s="$(tr -dc 'A-Z2-7' </dev/urandom 2>/dev/null | head -c 16)"
	if [ -n "$s" ] && [ "${#s}" -eq 16 ]; then
		printf '%s' "$s"
		return
	fi

	# Last-resort fallback. awk's srand() with no args seeds from the system
	# clock in busybox awk; without it we'd produce the SAME secret on every
	# boot, which is obviously catastrophic.
	awk 'BEGIN {
		srand()
		split("ABCDEFGHIJKLMNOPQRSTUVWXYZ234567", c, "")
		for (i = 1; i <= 16; i++)
			printf "%s", c[int(rand() * 32) + 1]
	}'
}

twofa_config_ok() {
	[ -f /etc/config/twofa ] || return 1
	uci -q show twofa >/dev/null 2>&1
}

# --- twofa: recreate broken/empty/whitespace-only config before any uci access ---
if ! twofa_config_ok; then
	write_twofa_config
	changed_twofa=1
fi

# Migrate legacy section type: config settings 'global'
typ="$(uci -q get twofa.global 2>/dev/null)" || typ=""
if [ "$typ" = "settings" ]; then
	uci -q set twofa.global=twofa
	changed_twofa=1
fi

if ! uci -q get twofa.global >/dev/null 2>&1; then
	uci -q set twofa.global=twofa
	uci -q set twofa.global.enabled='0'
	uci -q set twofa.global.issuer='OpenWrt'
	changed_twofa=1
fi

secret="$(uci -q get twofa.global.secret 2>/dev/null)" || secret=""
if [ -n "$secret" ] && [ -z "$(echo "$secret" | tr -d ' \t')" ]; then
	uci -q delete twofa.global.secret 2>/dev/null
	secret=""
	changed_twofa=1
fi

if [ -z "$secret" ]; then
	new_secret="$(generate_secret)"
	# Guard against generate_secret() somehow returning empty/short output -
	# committing secret="" would brick TOTP forever because the verify path
	# treats < 16 chars as "no secret configured".
	if [ -n "$new_secret" ] && [ "${#new_secret}" -ge 16 ]; then
		uci -q set twofa.global.secret="$new_secret"
		changed_twofa=1
	else
		logger -t "$name" "WARN: generate_secret() returned empty; leaving twofa.global.secret unset"
	fi
fi

[ "$changed_twofa" = "1" ] && uci -q commit twofa

# --- rpcd ACL for root (supports @login[0] and named login sections) ---
[ -f /etc/config/rpcd ] || exit 0

if uci -q show rpcd 2>/dev/null | grep -qF ".acl=${name}"; then
	exit 0
fi

for line in $(uci -q show rpcd 2>/dev/null | grep '=login$'); do
	sid="${line#rpcd.}"
	sid="${sid%%=*}"
	[ -n "$sid" ] || continue
	[ "$(uci -q get "rpcd.${sid}.username" 2>/dev/null)" = "root" ] || continue
	if ! uci -q show "rpcd.${sid}" 2>/dev/null | grep -qF ".acl=${name}"; then
		uci -q add_list "rpcd.${sid}.acl=${name}" 2>/dev/null && changed_rpcd=1
	fi
done

[ "$changed_rpcd" = "1" ] && uci -q commit rpcd

exit 0
