# root/etc/uci-defaults/10-luci-app-twofa
#!/bin/sh
# Initialize twofa config, migrate legacy in-UCI secret to /etc/twofa.secret,
# grant the rpcd ACL to every root login section, and kick the session-guard
# daemon. Idempotent: safe to re-run.

name="luci-app-twofa"
SECRET_FILE=/etc/twofa.secret
changed_twofa=0
changed_rpcd=0

write_twofa_config() {
	cat > /etc/config/twofa <<'EOF'
config twofa 'global'
	option enabled '0'
	option issuer 'OpenWrt'
EOF
}

twofa_config_ok() {
	[ -f /etc/config/twofa ] || return 1
	uci -q show twofa >/dev/null 2>&1
}

# --- twofa: recreate broken / missing / unparseable config ---
if ! twofa_config_ok; then
	write_twofa_config
	changed_twofa=1
fi

# Migrate legacy section type: `config settings 'global'`
typ="$(uci -q get twofa.global 2>/dev/null)" || typ=""
if [ "$typ" = "settings" ]; then
	uci -q set twofa.global=twofa
	changed_twofa=1
fi

if ! uci -q get twofa.global >/dev/null 2>&1; then
	uci -q set twofa.global=twofa
	uci -q set twofa.global.enabled='0'
	uci -q set twofa.global.issuer='OpenWrt'
	changed_twofa=1
fi

# --- v1.0-r16: secret moves out of UCI into /etc/twofa.secret (0600) ---
#
# UCI configs are world-readable (0644) by default and are commonly included
# in sysupgrade backups. Storing the TOTP secret there means anyone who
# obtains a backup can derive every future code. The secret now lives in
# /etc/twofa.secret with mode 0600 (root-only).
#
# Migration: copy the old uci-stored secret to the file (if absent), then
# unset it from UCI regardless so it never lingers.
legacy_secret="$(uci -q get twofa.global.secret 2>/dev/null)" || legacy_secret=""
if [ -n "$legacy_secret" ]; then
	if [ ! -s "$SECRET_FILE" ]; then
		legacy_secret_clean="$(printf '%s' "$legacy_secret" | tr -d ' \t\r\n')"
		if [ "${#legacy_secret_clean}" -ge 16 ]; then
			umask 077
			printf '%s\n' "$legacy_secret_clean" > "$SECRET_FILE"
			chmod 0600 "$SECRET_FILE"
			logger -t "$name" "migrated TOTP secret from uci to ${SECRET_FILE}"
		fi
	fi
	uci -q delete twofa.global.secret 2>/dev/null
	changed_twofa=1
fi

# Generate a fresh secret if the file is missing or too short.
need_gen=0
if [ ! -s "$SECRET_FILE" ]; then
	need_gen=1
elif [ "$(wc -c <"$SECRET_FILE" 2>/dev/null || echo 0)" -lt 16 ]; then
	need_gen=1
fi
if [ "$need_gen" = "1" ]; then
	if /usr/sbin/twofa-genkey >/dev/null 2>&1; then
		logger -t "$name" "generated new TOTP secret in ${SECRET_FILE}"
	else
		logger -t "$name" "WARN: twofa-genkey failed; 2FA will refuse to enable until a secret exists"
	fi
fi

# Defensive chmod even if file already existed (older installs may have
# created it with default umask).
[ -f "$SECRET_FILE" ] && chmod 0600 "$SECRET_FILE"

[ "$changed_twofa" = "1" ] && uci -q commit twofa

# --- rpcd ACL for root login sections ---
[ -f /etc/config/rpcd ] || exit 0

if uci -q show rpcd 2>/dev/null | grep -qF ".acl=${name}"; then
	# ACL already granted on at least one login section. Even so, walk all
	# root logins below so a freshly-added one still gets it.
	:
fi

for line in $(uci -q show rpcd 2>/dev/null | grep '=login$'); do
	sid="${line#rpcd.}"
	sid="${sid%%=*}"
	[ -n "$sid" ] || continue
	[ "$(uci -q get "rpcd.${sid}.username" 2>/dev/null)" = "root" ] || continue
	if ! uci -q show "rpcd.${sid}" 2>/dev/null | grep -qF ".acl=${name}"; then
		uci -q add_list "rpcd.${sid}.acl=${name}" 2>/dev/null && changed_rpcd=1
	fi
done

[ "$changed_rpcd" = "1" ] && uci -q commit rpcd

# Restart the guard daemon now that config + secret are in place. The
# package postinst also handles this, but uci-defaults can run on its own
# (factory reset, manual re-run) so we cover that case too.
if [ -x /etc/init.d/luci-app-twofa-guard ]; then
	/etc/init.d/luci-app-twofa-guard enable  >/dev/null 2>&1 || true
	/etc/init.d/luci-app-twofa-guard restart >/dev/null 2>&1 || true
fi

exit 0
